# src/base/etc/inittab
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
tty1::askfirst:/sbin/agetty 9600 tty1
tty2::askfirst:/sbin/agetty 9600 tty2
tty3::askfirst:/sbin/agetty 9600 tty3
tty4::askfirst:/sbin/agetty 9600 tty4
